# scripts/autoshield.service
[Unit]
Description=AutoShield Intrusion Prevention Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/autoshield/venv/bin/python -m src.main
Restart=on-failure
RestartSec=10
User=root
WorkingDirectory=/opt/autoshield

# Security settings
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
